# dwm/v2/systemd/user/easyeffects.service
[Unit]
Description=Easy Effects Service

[Service]
ExecStart=easyeffects --hide-window --service-mode
